# tools/livekit-publisher.service
[Unit]
Description=LiveKit Camera Publisher
After=wifi-qr-connect.service network-online.target
Wants=wifi-qr-connect.service network-online.target

[Service]
Type=simple
User=trc
WorkingDirectory=/home/trc
ExecStart=/home/trc/publish_livekit.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
